# deploy/systemd/fortress-ai-dashboard.service
[Unit]
Description=Fortress AI v2 Dashboard (Flask)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/fortress-ai
Environment=PYTHONPATH=/home/ubuntu/fortress-ai
Environment=FORTRESS_AI_DASHBOARD_HOST=0.0.0.0
Environment=FORTRESS_AI_DASHBOARD_PORT=8050
EnvironmentFile=-/home/ubuntu/fortress-ai/.env
ExecStart=/usr/bin/python3 /home/ubuntu/fortress-ai/dashboard/ai_command_center.py
Restart=on-failure
RestartSec=5
PrivateTmp=true

StandardOutput=append:/home/ubuntu/fortress-ai/logs/dashboard_systemd.log
StandardError=append:/home/ubuntu/fortress-ai/logs/dashboard_systemd.log

[Install]
WantedBy=multi-user.target
