# deploy/fortress-ai-operator-status.service
# Operator status monitor — 15-minute snapshots to data/operator_status/
# Install:
#   sudo cp deploy/fortress-ai-operator-status.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now fortress-ai-operator-status.service

[Unit]
Description=Fortress AI operator status monitor (15m loop)
After=network-online.target fortress-ai-skim-swarm.service fortress-ai-infra-swarm.service
Wants=fortress-ai-skim-swarm.service fortress-ai-infra-swarm.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/fortress-ai
Environment=PYTHONPATH=/home/ubuntu/fortress-ai
Environment=PATH=/home/ubuntu/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=TZ=America/New_York
Environment=FORTRESS_SYSTEM_TZ=America/New_York
EnvironmentFile=-/home/ubuntu/fortress-ai/.env
ExecStart=/bin/bash /home/ubuntu/fortress-ai/scripts/operator_status_loop.sh
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
