# deploy/fortress-ai-governance.timer
[Unit]
Description=Daily Fortress AI governance maintenance

[Timer]
OnCalendar=*-*-* 06:15:00
Persistent=true

[Install]
WantedBy=timers.target
